# config/systemd/mcp-agent-mail.service
[Unit]
Description=MCP Agent Mail Server
After=network.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__HOME__/mcp_agent_mail
Environment="PATH=__HOME__/.local/bin:__HOME__/.bun/bin:/usr/local/bin:/usr/bin:/bin"
Environment="HOME=__HOME__"
ExecStart=__HOME__/mcp_agent_mail/scripts/run_server_with_token.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
