# scripts/farmhand-cleanup.service
[Unit]
Description=Farmhand Stale Reservation Cleanup
Documentation=https://github.com/eaasxt/Farmhand
After=mcp-agent-mail.service

[Service]
Type=oneshot
ExecStart=/home/ubuntu/.local/bin/bd-cleanup --cron --stale-hours=2
User=ubuntu
Environment=HOME=/home/ubuntu
Environment=BEADS_DB=/home/ubuntu/.beads/beads.db

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=farmhand-cleanup
